# roles/walg/templates/walg-backup.service.j2
# walg-backup.service — managed by Ansible (walg role)
# Do not edit manually; changes will be overwritten.
[Unit]
Description=WAL-G full backup push to S3
After=postgresql.service
Requires=postgresql.service
ConditionPathExists={{ walg_config_dir }}/env

[Service]
Type=oneshot
User=postgres
EnvironmentFile={{ walg_config_dir }}/env
ExecStart={{ walg_install_path }} backup-push {{ pg_data_dir }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=walg-backup
TimeoutStartSec=3600
